# target/linux/imx/imx8/base-files/etc/uci-defaults/20_configure-wwan
#!/bin/sh

. /lib/functions.sh

DEFAULT_APN="fast.t-mobile.com"

uci -q batch << EOI
set network.wan.proto='mbim'
set network.wan.device='/dev/cdc-wdm0'
set network.wan.apn=$DEFAULT_APN
set network.wan.auth='none'
set network.wan.pdptype='ipv4v6'
delete network.wan6
EOI

uci -q commit network