# target/linux/imx/imx8/base-files/etc/uci-defaults/21_configure-uap
wifi down

uci set wireless.radio0.country='US'
uci set wireless.radio0.htmode='HT40'
uci set wireless.radio0.disabled='0'
uci commit wireless

wifi up