# target/linux/imx/imx8/base-files/etc/uci-defaults/30_configure-gpsd
uci set gpsd.core=gpsd
uci set gpsd.core.device='/dev/ttymxc0'
uci set gpsd.core.port='2947'
uci set gpsd.core.listen_globally='0'
uci set gpsd.core.disabled='0'

uci commit gpsd