# target/linux/imx/imx8/base-files/etc/uci-defaults/30_configure-gpsd
#!/bin/sh

. /lib/functions.sh

GPS_DEVICE="/dev/ttymxc0"

uci -q batch << EOI
set gpsd.core.device=$GPS_DEIVCE
set gpsd.core.listen_globally='0'
set gpsd.core.enabled='1'
EOI

uci -q commit gpsd

/etc/init.d/gpsd enable
/etc/init.d/gpsd start