# target/linux/imx/imx8/base-files/etc/uci-defaults/50_configure-default-access
#!/bin/sh

. /lib/functions.sh

# Allow SSH access via LAN devices
uci set dropbear.@dropbear[0].Interface='lan'
uci -q commit dropbear

# HTTPS redirect for LuCI
uci set uhttpd.main.redirect_https='1'
uci -q commit uhttpd
