# target/linux/imx/imx8/base-files/etc/uci-defaults/55_configure-gpsd
#!/bin/sh

. /lib/functions.sh

GPS_DEVICE="/dev/ttymxc0"

uci -q batch << EOI
set gpsd.core.device=$GPS_DEVICE
set gpsd.core.enabled='1'
EOI

uci -q commit gpsd